# supervisord.conf
; TVAiPlatform — Supervisor configuration
; Usage: supervisord -c supervisord.conf

[unix_http_server]
file=./data/logs/supervisor.sock

[supervisord]
logfile=./data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=./data/logs/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/data/logs/supervisor.sock

; ---- Services ----

[program:tvai-api]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8000
directory=%(ENV_PROJECT_ROOT)s/api
autostart=true
autorestart=true
startsecs=5
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/api.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/api-error.log
environment=PYTHONUNBUFFERED="1"

[program:tvai-web]
command=pnpm dev --port 3100
directory=%(ENV_PROJECT_ROOT)s/web
autostart=true
autorestart=true
startsecs=5
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/web.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/web-error.log

[program:chattts]
command=python chattts_server.py
directory=%(ENV_PROJECT_ROOT)s/tts
autostart=true
autorestart=true
startsecs=10
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/chattts.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/chattts-error.log
environment=PYTHONUNBUFFERED="1"

[program:redis]
command=redis-server --daemonize no --dir %(ENV_PROJECT_ROOT)s/data/redis
autostart=true
autorestart=true
startsecs=3
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/redis.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/redis-error.log

[program:postgresql]
command=postgres -D /var/lib/postgresql/17/main
autostart=true
autorestart=true
startsecs=5
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/postgres.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/postgres-error.log

[program:temp-cleaner]
command=bash %(ENV_PROJECT_ROOT)s/scripts/cleanup_temp.sh
directory=%(ENV_PROJECT_ROOT)s
autostart=true
startretries=0
startsecs=0
exitcodes=0
stdout_logfile=%(ENV_PROJECT_ROOT)s/data/logs/cleanup.log
stderr_logfile=%(ENV_PROJECT_ROOT)s/data/logs/cleanup-error.log

[group:tvai]
programs=tvai-api,tvai-web,chattts,redis,postgresql
